# astrobase_services_pip/astrobase_services/supervisor_astrobase_services.ini
[group:atdb_services_data]
programs=data_monitor, cleanup

[group:atdb_services_ingest]
programs=start_ingest, start_ingest_priority, start_ingest_small

[program:data_monitor]
command=/bin/bash -c 'cd /home/apertif/atdb_client; source env2/bin/activate; source /opt/apertif/apertifinit.sh; source /opt/gcc-6.3.0/gcc.rc; atdb_service -o data_monitor --obs_mode_filter imaging --host_filter wcudata1 --interval 30 --atdb_host prod -v'
environment=HOME=/home/apertif
user=apertif
stopsignal=KILL
stopasgroup=true ; also kills the subprocesses
stdout_logfile=/home/apertif/atdb_client/atdb_service_data_monitor.log
redirect_stderr=true
stderr_logfile=NONE

[program:start_ingest]
command=/bin/bash -c 'cd /home/apertif/atdb_client; source env2/bin/activate; source /opt/apertif/apertifinit.sh; source /opt/gcc-6.3.0/gcc.rc; atdb_service -o start_ingest --obs_mode_filter imaging --status_filter valid --interval 30 --atdb_host prod --alta_host prod -v'
environment=HOME=/home/apertif,LD_LIBRARY_PATH="/home/apertif/alta-client/.env/lib:$LD_LIBRARY_PATH"
user=apertif
stopsignal=KILL
stopasgroup=true ; also kills the subprocesses
stdout_logfile=/home/apertif/atdb_client/atdb_service_start_ingest.log
redirect_stderr=true
stderr_logfile=NONE
autostart=false

[program:start_ingest_small]
command=/bin/bash -c 'cd /home/apertif/atdb_client; source env2/bin/activate; source /opt/apertif/apertifinit.sh; source /opt/gcc-6.3.0/gcc.rc; atdb_service -o start_ingest --obs_mode_filter imaging --status_filter valid_small --interval 60 --atdb_host prod --alta_host prod -v'
environment=HOME=/home/apertif,LD_LIBRARY_PATH="/home/apertif/alta-client/.env/lib:$LD_LIBRARY_PATH"
user=apertif
stopsignal=KILL
stopasgroup=true ; also kills the subprocesses
stdout_logfile=/home/apertif/atdb_client/atdb_service_start_ingest_small.log
redirect_stderr=true
stderr_logfile=NONE
autostart=false

[program:start_ingest_priority]
command=/bin/bash -c 'cd /home/apertif/atdb_client; source env2/bin/activate; source /opt/apertif/apertifinit.sh; source /opt/gcc-6.3.0/gcc.rc; atdb_service -o start_ingest --obs_mode_filter imaging --status_filter valid_priority --interval 60 --atdb_host prod --alta_host prod -v'
environment=HOME=/home/apertif,LD_LIBRARY_PATH="/home/apertif/alta-client/.env/lib:$LD_LIBRARY_PATH"
user=apertif
stopsignal=KILL
stopasgroup=true ; also kills the subprocesses
stdout_logfile=/home/apertif/atdb_client/atdb_service_start_ingest_priority.log
redirect_stderr=true
stderr_logfile=NONE
autostart=false

[program:cleanup]
command=/bin/bash -c 'cd /home/apertif/atdb_client; source env2/bin/activate; source /opt/apertif/apertifinit.sh; source /opt/gcc-6.3.0/gcc.rc; atdb_service -o cleanup --obs_mode_filter imaging --interval 60 --atdb_host prod'
environment=HOME=/home/apertif
user=apertif
stopsignal=KILL
stopasgroup=true ; also kills the subprocesses
stdout_logfile=/home/apertif/atdb_client/atdb_service_cleanup.log
redirect_stderr=true
stderr_logfile=NONE
